# service_php/service_php.sv.conf
[program:service_php]
command=php /usr/src/service_php/index.php
redirect_stderr=true
